# ubuntu/zram.conf
start on started lxc-android-config

task

script
zramctl -s 800M /dev/zram0 
mkswap /dev/zram0
swapon /dev/zram0

end script
